# files/paperless_consumerd.in
#!/bin/sh

# PROVIDE: paperless_consumerd
# REQUIRE: DAEMON LOGIN redis
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable paperless_consumerd:
# paperless_consumerd_enable="YES"
# Or run:
# sysrc paperless_consumerd_enable="YES"

. /etc/rc.subr

name="paperless_consumerd"
rcvar="${name}_enable"
procname="daemon"

pidfile="/var/run/${name}.pid"

command="/usr/sbin/daemon"
command_args="-f -P ${pidfile} -u %%USER%% %%OPTDIR%%/bin/paperless_consumer"

load_rc_config $name
run_rc_command "$1"
