# files/paperless_schedulerd.in
#!/bin/sh

# PROVIDE: paperless_schedulerd
# REQUIRE: DAEMON LOGIN redis
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable paperless_schedulerd:
# paperless_schedulerd_enable="YES"
# Or run:
# sysrc paperless_schedulerd_enable="YES"

. /etc/rc.subr

name="paperless_schedulerd"
rcvar="${name}_enable"
procname="daemon"

pidfile="/var/run/${name}.pid"

command="/usr/sbin/daemon"
command_args="-f -P ${pidfile} -u %%USER%% %%OPTDIR%%/bin/paperless_scheduler"

load_rc_config $name
run_rc_command "$1"
